# bridges/FilesystemRouterBridge/systemd/filesystem-router.service
[Unit]
# FilesystemRouterBridge — Periodic file routing between Bridgit bridges.
#
# Routes packages from bridge outboxes to bridge inboxes based on the
# route definitions in registry/routes.toml. Fires every 30 seconds
# to ensure timely delivery without excessive filesystem polling.
Description=Bridgit Filesystem Router Bridge
After=network.target

[Service]
Type=oneshot
ExecStart=/home/craigpar/Code/bridgit-sync-engine/bridges/FilesystemRouterBridge/bin/route.sh
WorkingDirectory=/home/craigpar/Code/bridgit-sync-engine
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
